# systemd/cm108d.service
[Unit]
Description=CM108/CM119 real-time audio server
After=local-fs.target

[Service]
Type=simple
ExecStart=/usr/sbin/cm108d --socket /run/cm108d/cm108d.sock --log info
RuntimeDirectory=cm108d
RuntimeDirectoryMode=0750

# Real-time privileges
AmbientCapabilities=CAP_SYS_NICE CAP_IPC_LOCK
LimitRTPRIO=95
LimitMEMLOCK=infinity

Restart=on-failure
RestartSec=2s

# Hardening
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
